# softcursor.service
# softcursor — VMware ESXi 6.x software cursor overlay
# Built WITH CSK by NoHatHacker.com

[Unit]
Description=Software cursor overlay for VMware ESXi 6.x guests
Documentation=https://github.com/mr1b31r0/Kali_ESXi_mouseFixer
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStartPre=/bin/sh -c 'test -n "$DISPLAY" || export DISPLAY=:0'
ExecStart=%h/.local/bin/softcursor
Restart=on-failure
RestartSec=3
# Detect DISPLAY from running X server socket if not inherited
Environment=DISPLAY=:0

[Install]
WantedBy=graphical-session.target
WantedBy=default.target
